# packaging/debian/DEBIAN/postinst
#!/bin/bash
set -e

AUDYN_USER="audyn"
AUDYN_GROUP="audyn"
AUDYN_HOME="/var/lib/audyn"
AUDYN_ARCHIVE="/var/lib/audyn/archive"
AUDYN_CONFIG="/etc/audyn"
AUDYN_LOG="/var/log/audyn"
AUDYN_OPT="/opt/audyn"

case "$1" in
    configure)
        # Create system user and group
        if ! getent group "$AUDYN_GROUP" > /dev/null 2>&1; then
            groupadd --system "$AUDYN_GROUP"
            echo "Created system group: $AUDYN_GROUP"
        fi

        if ! getent passwd "$AUDYN_USER" > /dev/null 2>&1; then
            useradd --system \
                --gid "$AUDYN_GROUP" \
                --home-dir "$AUDYN_HOME" \
                --shell /usr/sbin/nologin \
                --comment "Audyn Audio Capture Service" \
                "$AUDYN_USER"
            echo "Created system user: $AUDYN_USER"
        fi

        # Create directories
        mkdir -p "$AUDYN_HOME"
        mkdir -p "$AUDYN_ARCHIVE"
        mkdir -p "$AUDYN_CONFIG"
        mkdir -p "$AUDYN_LOG"

        # Set ownership and permissions
        chown -R "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_HOME"
        chown -R "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_ARCHIVE"
        chown -R "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_LOG"
        chown root:$AUDYN_GROUP "$AUDYN_CONFIG"
        chmod 750 "$AUDYN_HOME"
        chmod 750 "$AUDYN_ARCHIVE"
        chmod 750 "$AUDYN_LOG"
        chmod 775 "$AUDYN_CONFIG"

        # Create default configuration if not exists
        if [ ! -f "$AUDYN_CONFIG/global.json" ]; then
            cat > "$AUDYN_CONFIG/global.json" << 'EOFCONFIG'
{
  "source_type": "aes67",
  "multicast_addr": "239.69.1.1",
  "port": 5004,
  "sample_rate": 48000,
  "channels": 2,
  "format": "wav",
  "bitrate": 128000,
  "archive_root": "/var/lib/audyn/archive",
  "archive_layout": "dailydir",
  "archive_period": 3600,
  "archive_clock": "localtime",
  "ptp_interface": null,
  "aes67_interface": null
}
EOFCONFIG
            chown "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_CONFIG/global.json"
            chmod 640 "$AUDYN_CONFIG/global.json"
        fi

        # Create default system configuration
        if [ ! -f "$AUDYN_CONFIG/system.json" ]; then
            cat > "$AUDYN_CONFIG/system.json" << 'EOFSYSTEM'
{
  "hostname": null,
  "timezone": "UTC",
  "ntp_servers": ["pool.ntp.org"]
}
EOFSYSTEM
            chown "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_CONFIG/system.json"
            chmod 640 "$AUDYN_CONFIG/system.json"
        fi

        # Create default SSL configuration
        if [ ! -f "$AUDYN_CONFIG/ssl.json" ]; then
            cat > "$AUDYN_CONFIG/ssl.json" << 'EOFSSL'
{
  "enabled": false,
  "domain": null,
  "email": null,
  "auto_renew": true
}
EOFSSL
            chown "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_CONFIG/ssl.json"
            chmod 640 "$AUDYN_CONFIG/ssl.json"
        fi

        # Create default auth configuration with breakglass password
        if [ ! -f "$AUDYN_CONFIG/auth.json" ]; then
            cat > "$AUDYN_CONFIG/auth.json" << 'EOFAUTH'
{
  "entra_tenant_id": null,
  "entra_client_id": null,
  "entra_redirect_uri": null,
  "breakglass_password_hash": "$2b$12$G43XJz5S3qJodySjUgu0LOvh4Ns5xxZ2.me7Wspv9z4QmL7KOxpte"
}
EOFAUTH
            chown "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_CONFIG/auth.json"
            chmod 640 "$AUDYN_CONFIG/auth.json"
            echo ""
            echo "  DEFAULT BREAKGLASS PASSWORD: audyn"
            echo "  Please change this immediately after login!"
            echo ""
        fi

        # Set up Python virtual environment
        if [ ! -d "$AUDYN_OPT/backend/venv" ]; then
            echo "Setting up Python virtual environment..."
            python3 -m venv "$AUDYN_OPT/backend/venv"
            "$AUDYN_OPT/backend/venv/bin/pip" install --upgrade pip
            "$AUDYN_OPT/backend/venv/bin/pip" install -r "$AUDYN_OPT/backend/requirements.txt"
        fi

        # Set permissions on opt directory
        chown -R "$AUDYN_USER:$AUDYN_GROUP" "$AUDYN_OPT"

        # Configure Samba share for archive access
        SAMBA_CONF="/etc/samba/smb.conf"
        if ! grep -q "\[audyn\]" "$SAMBA_CONF" 2>/dev/null; then
            echo "Configuring Samba share for archive access..."

            # Backup original config
            if [ -f "$SAMBA_CONF" ]; then
                cp "$SAMBA_CONF" "${SAMBA_CONF}.bak"
            fi

            # Add audyn share configuration
            cat >> "$SAMBA_CONF" << 'EOFSAMBA'

# Audyn Archive Share
[audyn]
   comment = Audyn Audio Archive
   path = /var/lib/audyn/archive
   browseable = yes
   read only = yes
   guest ok = no
   valid users = audyn
EOFSAMBA

            # Ensure SMB3 minimum protocol in global section (security)
            if ! grep -q "server min protocol" "$SAMBA_CONF"; then
                sed -i '/^\[global\]/a\   server min protocol = SMB3' "$SAMBA_CONF"
            fi

            echo "Samba share configuration added"
        fi

        # Create Samba user (separate from system user)
        # The system user was created with nologin shell, but Samba auth is separate
        if ! pdbedit -L 2>/dev/null | grep -q "^audyn:"; then
            echo "Creating Samba user: audyn"
            # Set password to 'audyn' non-interactively
            (echo "audyn"; echo "audyn") | smbpasswd -a -s audyn
            smbpasswd -e audyn
            echo ""
            echo "  DEFAULT SAMBA PASSWORD: audyn"
            echo "  Change with: sudo smbpasswd audyn"
            echo ""
        fi

        # Enable and start Samba services
        systemctl enable smbd.service nmbd.service
        systemctl restart smbd.service nmbd.service

        # Enable nginx site
        if [ -f /etc/nginx/sites-available/audyn ]; then
            ln -sf /etc/nginx/sites-available/audyn /etc/nginx/sites-enabled/audyn
            # Remove default site if exists
            rm -f /etc/nginx/sites-enabled/default
        fi

        # Reload systemd and enable services
        systemctl daemon-reload
        systemctl enable audyn-web.service
        systemctl enable nginx.service

        # Start services
        systemctl restart nginx.service
        systemctl start audyn-web.service

        echo ""
        echo "========================================"
        echo "  Audyn installation complete!"
        echo "========================================"
        echo ""
        echo "Web interface: http://<management-ip>/"
        echo ""
        echo "Login with breakglass password: audyn"
        echo "CHANGE THIS PASSWORD IMMEDIATELY!"
        echo ""
        echo "SMB Share: \\\\<management-ip>\\audyn"
        echo "SMB User: audyn  Password: audyn"
        echo "(Change with: sudo smbpasswd audyn)"
        echo ""
        echo "Default archive: $AUDYN_ARCHIVE"
        echo "Configuration: $AUDYN_CONFIG"
        echo "Logs: $AUDYN_LOG"
        echo ""
        ;;

    abort-upgrade|abort-remove|abort-deconfigure)
        ;;

    *)
        echo "postinst called with unknown argument: $1" >&2
        exit 1
        ;;
esac

exit 0
